# resources/autostart/desktop-refresh.desktop
[Desktop Entry]
Type=Application
Name=Desktop Refresh
Description=
Exec=/home/breeze256/.local/share/scripts/desktop-refresh.sh
