# resources/autostart/nemo-desktop.desktop
[Desktop Entry]
Type=Application
Name=Nemo Desktop
Description=
Exec=nemo-desktop
